# supervisord.ini
[program:leadrouter_subscriber]
command=/usr/bin/leadrouter subscriber --beanstalkd-host=%(ENV_BEANSTALK_HOST)s --beanstalkd-port=%(ENV_BEANSTALK_PORT)s --beanstalkd-tube=%(ENV_BEANSTALK_TUBE)s --loglevel=%(ENV_LOGLEVEL)s
numprocs=%(ENV_WORKERS)s
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
